# infrastructure/backend/radio-api.service
[Unit]
Description=Radio API Service
After=network.target

[Service]
User=root
WorkingDirectory=/app
# CRITICAL: Added --proxy-headers and --forwarded-allow-ips for Nginx Websocket support
ExecStart=/opt/radio-api/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips "*"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
